# docker_kippo/docker-armhf-kippo-restart.timer
[Unit]
Description=Docker Container armhf-kippo restarter

[Timer]
OnBootSec=1h
OnUnitActiveSec=2d

[Install]
WantedBy=multi-user.target
